# Jetson Orin Nano Carrier Board/SG6C-ORNX-G2-F/JetPack5.1.2/SG6C-ORNX-G2-F_OrinNX_OrinNano_ISX031_MAX96717Fx6_JP5.1.2_L4TR35.4.1_1920x1536/ko/load-pps-gpio.service
[Unit]
Description=Load pps-gpio Kernel Module

[Service]
ExecStart=/sbin/insmod /lib/modules/5.10.120-tegra/kernel/drivers/pps/clients/pps-gpio.ko
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

